# .config/systemd/user/usermount.service
[Unit]
Description=Mount removable devices.

[Service]
ExecStart=/usr/bin/usermount

[Install]
WantedBy=default.target
