# .config/systemd/user/xscreensaver.service
[Unit]
Description=XScreensaver daemon

[Service]
ExecStart=/usr/bin/xscreensaver -nosplash

[Install]
WantedBy=i3-session.target
